# fhem/contrib/DEBIAN/postinst
#!/bin/bash
# $Id$

set -e
#set -x



if ! getent passwd fhem >/dev/null; then
    useradd --system --home /opt/fhem --gid dialout --shell /bin/false fhem
fi

chown -R fhem:dialout /opt/fhem

# remove repository entry after first installation
sed -i /debian.fhem.de/d /etc/apt/sources.list

# set up of autostart
# Sysvinit; upstart no longer supported

#mkdir -p /etc/init.d
cp /opt/fhem/contrib/init-scripts/fhem.3 /etc/init.d/fhem
chmod ugo+x /etc/init.d/fhem
update-rc.d fhem defaults

if test -f /etc/init.d/fhem; then
    /etc/init.d/fhem start noaptmark
fi
exit 0
