# fhem/contrib/DEBIAN/prerm
#!/bin/sh

set -e

# stop fhem if running
invoke-rc.d fhem stop

# remove fhem autostart
update-rc.d fhem remove
