# fhem/contrib/DEBIAN/prerm
#!/bin/sh

#set -x
set -e

apt-mark unhold fhem

# stop fhem if running
if test -f /etc/init.d/fhem; then
    /etc/init.d/fhem stop noaptmark
fi

# remove fhem autostart
update-rc.d -f fhem remove
