# fhem/contrib/DEBIAN/prerm
#!/bin/sh

set -e

apt-mark unhold fhem

# stop fhem if running
if test -f /etc/init.d/fhem; then
    invoke-rc.d fhem stop
fi

# remove fhem autostart
update-rc.d fhem remove
